# scripts/udev/10-clj.rules

KERNEL=="ttyUSB*", SUBSYSTEM=="tty", RUN+="/opt/ablebit/clj/bin/clj-ppp.sh %k %n %p %b $driver %M %m $name %S"
