# scripts/udev/10-clj.rules

KERNEL=="ttyUSB*", SUBSYSTEM=="tty", RUN+="clj-ppp.sh %k %n %p %b $driver %M %m $name %S"

